# packaging/systemd/pfs-maint@.service
[Unit]
Description=PolicyFS maint (mount %i)

[Service]
Type=oneshot
ExecStart=/usr/bin/pfs maint %i
SuccessExitStatus=75 3
