# packaging/systemd/pfs-move@.service
[Unit]
Description=PolicyFS mover (mount %i)

[Service]
Type=oneshot
ExecStart=/usr/bin/pfs move %i
SuccessExitStatus=75 3
