# packaging/systemd/pfs-index@.service
[Unit]
Description=PolicyFS index (mount %i)

[Service]
Type=oneshot
ExecStart=/usr/bin/pfs index %i
SuccessExitStatus=75 3
