# packaging/systemd/pfs-prune@.service
[Unit]
Description=PolicyFS prune (mount %i)

[Service]
Type=oneshot
ExecStart=/usr/bin/pfs prune %i
SuccessExitStatus=75 3
