# packaging/systemd/pfs@.service
[Unit]
Description=PolicyFS daemon (mount %i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
RuntimeDirectory=pfs/%i
RuntimeDirectoryMode=0755
StateDirectory=pfs/%i
StateDirectoryMode=0755
LogsDirectory=pfs
LogsDirectoryMode=0755
ExecStart=/usr/bin/pfs mount %i
ExecStop=/usr/bin/pfs unmount %i
TimeoutStopSec=10s
Restart=on-failure
RestartSec=2s

# Treat "lock already held" as a clean failure (avoid spamming logs)
SuccessExitStatus=75

[Install]
WantedBy=multi-user.target
